# app-web-1.conf
start on starting app-web
stop on stopping app-web
respawn

env PORT=5000

setuid app

chdir /Users/ytyeoh/Desktop/pair_bnb

exec bundle exec puma -p $PORT
